# 99-keyboard.rules
ACTION=="add", SUBSYSTEM=="input", RUN+="/bin/su -l krvkir -c '/home/krvkir/bin/keyboard-setup.sh'"

